# systemd/user/gnome-terminal-theme.service
[Unit]
Description=Apply GNOME Terminal time-based theme
After=graphical-session.target dbus.socket
PartOf=graphical-session.target
ConditionUser=!@system

[Service]
Type=oneshot
ExecStart=%h/scripts/terminal-theme.sh
StandardOutput=append:%h/scripts/terminal-theme.log
StandardError=append:%h/scripts/terminal-theme.log

[Install]
WantedBy=default.target
